# scripts/service/kresd/run
#!/bin/sh

set -eu

cd ~/
exec 2>&1
exec kresd --forks=1 --config=/etc/knot-resolver/kresd.conf
