# scripts/service/kres-cache-gc/run
#!/bin/sh

set -eu

cd "${HOME:?}"
exec 2>&1
exec kres-cache-gc -c "${HOME:?}" -d 1000
