# scripts/service/supercronic/run
#!/bin/sh

set -eu

cd "${KRESD_CACHE_DIR:?}"

exec 2>&1
exec supercronic /etc/crontab
