# scripts/service/supercronic/run
#!/bin/sh

set -eu

cd ~/
exec 2>&1
exec supercronic /etc/crontab
